# s6-overlay/s6-rc.d/cloudcli/run
#!/command/with-contenv sh
# shellcheck shell=sh
cd /workspace || exit 1
export HOME=/home/claude
export WORKSPACES_ROOT=/workspace
export NODE_OPTIONS="${NODE_OPTIONS:+$NODE_OPTIONS }--no-deprecation"
if [ "$(id -u)" = "0" ]; then
  exec s6-setuidgid claude cloudcli --port 3001
fi
exec cloudcli --port 3001
