# s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv sh
# shellcheck shell=sh
exec /usr/sbin/sshd -D -e -f /etc/ssh/sshd_config_holyclaude
